# files/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:iventoy]
directory=/iventoy
command=./lib/iventoy
environment=IVENTOY_API_ALL=1,IVENTOY_AUTO_RUN=1,LD_LIBRARY_PATH=/iventoy/lib/lin64
